# systemd/user/llama-server.service
[Unit]
Description=llama.cpp API Server (Qwen2.5-Coder-7B-Instruct on CPU)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/sstelfox/workspace/scratch/llama.cpp
ExecStart=/home/sstelfox/workspace/scratch/llama.cpp/build/bin/llama-server \
  -hfr Qwen/Qwen2.5-Coder-7B-Instruct-GGUF \
  -hff qwen2.5-coder-7b-instruct-q4_k_m.gguf \
  --port 17200 \
  --host 0.0.0.0 \
  -ngl 0 \
  --sleep-idle-seconds 300 \
  -t 20 \
  -tb 20 \
  -c 131072 \
  -b 4096 \
  -ub 2048 \
  -np 1 \
  --mlock \
  --cont-batching \
  -fa on \
  --prio 2 \
  --poll 100
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
